# deploy/anki-bot.service
[Unit]
Description=Anki Daily Bot — Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/anki-daily-bot
EnvironmentFile=/root/anki-daily-bot/.env
ExecStart=/root/anki-daily-bot/.venv/bin/python -m src.bot_runner
Restart=on-failure
RestartSec=10
StandardOutput=append:/root/anki-daily-bot/logs/bot.log
StandardError=append:/root/anki-daily-bot/logs/bot.log

[Install]
WantedBy=multi-user.target
