# deploy/anki-headless.service
[Unit]
Description=Anki (headless via Xvfb)
After=network-online.target xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=anki
Environment=DISPLAY=:99
Environment=HOME=/home/anki
ExecStart=/usr/local/bin/anki
ExecStartPost=/bin/bash -c 'sleep 20 && DISPLAY=:99 xdotool key Escape Return Escape Return'
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
