# deploy/anki-bot.service
[Unit]
Description=Anki Daily Bot — Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/anki-daily-bot
EnvironmentFile=/home/YOUR_USER/anki-daily-bot/.env
ExecStart=/home/YOUR_USER/anki-daily-bot/.venv/bin/python -m src.bot_runner
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/YOUR_USER/anki-daily-bot/logs/bot.log
StandardError=append:/home/YOUR_USER/anki-daily-bot/logs/bot.log

[Install]
WantedBy=multi-user.target
